# systemd/system/netplan-daily.service
[Unit]
Description=Daily Netplan Apply

[Service]
Type=oneshot
ExecStart=/usr/sbin/netplan apply
